# wavealert360-startup.service
[Unit]
Description=WaveAlert360 Beach Safety Alert System - Startup
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/WaveAlert360
ExecStart=/home/pi/WaveAlert360/restart_system.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="PATH=/home/pi/WaveAlert360/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Hardening options
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
